# autostart/Keymapp.desktop
[Desktop Entry]
Type=Application
Name=Keymapp
Exec="/usr/bin/keymapp"
X-GNOME-Autostart-enabled=true
